# RASPBERRY_sm/autostart/ini_main_py.desktop
[Desktop Entry]
Exec=python3 /home/pi/Desktop/PY_hab_9_16/main_hab_9_16.py
