# RASPBERRY_sm/autostart/ini_main_java.desktop
[Desktop Entry]
Exec=java -jar /home/pi/Desktop/dist/SisMotel_RPI.jar